# cron.d/logster
# execute logster every  minute
* * * * * root /usr/bin/logster
# done
