# cron.d/logster
# execute logster every  minute
* * * * * root /usr/bin/logster >/dev/null 2>/dev/null
# done
